# net-misc/sqm-scripts/files/sqm
#!/sbin/openrc-run
# Copyright 1999-2016 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

command="/usr/bin/sqm"

start() {
	if [ "${SVCNAME}" = "sqm" ]; then
		eerror "You are not supposed to run this script directly. Create a symlink"
		eerror
		eerror "  cd /etc/init.d"
		eerror "  ln -s sqm sqm.eth0"
		eerror
		return 1
	fi

	IFACE=${SVCNAME#*.}
	ebegin "Starting sqm on $IFACE"
	${command} start ${IFACE}
	eend $? "Failed to start sqm on $IFACE"
}

stop() {
	IFACE=${SVCNAME#*.}
	ebegin "Stopping sqm on $IFACE"
	${command} stop $IFACE
	eend $? "Failed to stop sqm on $IFACE"
}
